# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
command=/app/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

[program:worker]
command=/app/.venv/bin/celery -A app.worker worker --loglevel=info
; autostart is off until a later slice wires REDIS_URL through to Cloud Run and adds
; real pipeline tasks — starting it now would crash-loop (REDIS_URL isn't deployed yet)
; while /health stays green, masking the failure. Flip to true once REDIS_URL is set
; via --set-secrets/--set-env-vars in the deploy workflows.
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
