# supervisord.conf
[supervisord]
nodaemon=true

[program:app]
; --forwarded-allow-ips='*' trusts the X-Forwarded-Proto header from whatever peer connects to
; the container. Cloud Run terminates TLS at its own front end and always proxies to the
; container over a private, single-hop connection - the container is never reachable except
; through that proxy - so this is safe here. Without it, uvicorn's default (only trust
; 127.0.0.1) ignores X-Forwarded-Proto on Cloud Run's real proxy IP, and request.base_url
; (used to build the forgot-password reset link in app/auth/users.py) would report http://
; instead of https:// in prod/QA.
command=/app/.venv/bin/uvicorn app.main:app --host 0.0.0.0 --port 8000 --proxy-headers --forwarded-allow-ips='*'
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

; R13 (#168): the Celery worker (`app.worker`, `[program:worker]`) that used to run alongside the
; app here was removed along with the rest of the Host's dead event-extraction code - pipeline
; processing is now entirely event-creator's own concern (Cloud Tasks-based, per ADR-0001), not
; something the Host's container runs.
